# launchd/ai.yourgent.opencode.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.{{AGENT_NAME}}.opencode</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{OPENCODE_PATH}}</string>
        <string>serve</string>
        <string>--port</string>
        <string>8642</string>
        <string>--hostname</string>
        <string>0.0.0.0</string>
        <string>--print-logs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{DREWGENT_DIR}}</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{PATH}}</string>
        <key>OPENCODE_API_KEY</key>
        <string>{{OPENCODE_API_KEY}}</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>ThrottleInterval</key>
        <integer>10</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>{{DREWGENT_DIR}}/logs/opencode.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>{{DREWGENT_DIR}}/logs/opencode.stderr.log</string>
</dict>
</plist>
